# .config/systemd/user/keybase.gui.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/lib/opt/keybase/Keybase
